# patchbay-server/patchbay-serve.service
[Unit]
Description=patchbay-serve - patchbay run results server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=patchbay
Group=patchbay
ExecStart=/usr/local/bin/patchbay-serve \
    --accept-push \
    --data-dir /var/lib/patchbay-serve \
    --http-bind 0.0.0.0:80 \
    --https-bind 0.0.0.0:443 \
    --acme-domain patchbay.example.com \
    --acme-email you@example.com \
    --retention 10GB
Environment=PATCHBAY_API_KEY=changeme
AmbientCapabilities=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/patchbay-serve

[Install]
WantedBy=multi-user.target
